# collatzconjecture.service
[Unit]
Description=Collatz Conjecture
[Service]
Type=forking
ExecStart=/var/www/Collatz-Conjecture/src/run
[Install]
WantedBy=multi-user.target
